# files/usr/lib/systemd/system/usr-share-sddm-themes.mount
[Unit]
Description=OverlayFS for SDDM Themes
After=systemd-tmpfiles-setup.service
Before=sddm.service

[Mount]
What=overlay
Where=/usr/share/sddm/themes
Type=overlay
# Adjusted to match the Bazzite directory names
Options=lowerdir=/usr/share/sddm/themes,upperdir=/var/sddm_themes/themes,workdir=/var/sddm_themes/themes.work,index=off,metacopy=off

[Install]
WantedBy=multi-user.target
